# data/indicator-secret-agent.conf.in
description "Secret agent to repsond to Network Manager password requests"

start on starting indicator-network
stop on stopping indicator-network

respawn

exec @CMAKE_INSTALL_FULL_LIBEXECDIR@/indicator-secret-agent
