# data/indicator-secret-agent.conf.in
description "Secret agent to repsond to Network Manager password requests"

start on starting indicator-network
stop on stopping indicator-network

respawn

# Currently conflicts with the secret agent in Unity8
# so we're disabling until we can migrate fully.

# exec @CMAKE_INSTALL_FULL_LIBEXECDIR@/indicator-secret-agent

exec sleep 1
