# data/indicator-network-secret-agent.conf.in
description "Indicator Network Secret Agent"

start on starting indicator-network
stop on stopping indicator-network

respawn

exec @AGENT_SERVICE@
