# data/indicator-network.conf.in
description "Indicator Network Backend"

start on indicators-loaded

respawn

exec @NETWORK_SERVICE@
