# data/indicator-network-secret-agent.service.in
[Unit]
Description=Indicator Network Secret Agent
PartOf=indicator-network.service

[Service]
ExecStart=@AGENT_SERVICE@
Restart=on-failure
