# data/indicator-network.conf.in
description "Indicator Network Backend"

start on indicator-services-start
stop on desktop-end or indicator-services-end

env G_MESSAGES_DEBUG=all
export G_MESSAGES_DEBUG

respawn
respawn limit 2 10

exec @NETWORK_SERVICE@
