# data/indicator-network.service.in
[Unit]
Description=Indicator Network Backend
PartOf=graphical-session.target
After=indicators-pre.target unity8.service
Requisite=unity8.service
Wants=indicator-network-secret-agent.service

[Service]
ExecStart=@NETWORK_SERVICE@
Restart=on-failure
